# thesillyhome-addon/rootfs/etc/services.d/thesillyhome/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
echo "Hello world!"